# deepfence_agent/supervisord-cluster.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:deepfenced]
environment=DF_CLUSTER_AGENT=true
command=/bin/deepfenced
numprocs=1
autostart=true
autorestart=true
startsecs=0
stopwaitsecs=30
killasgroup=true
priority=50
